# deploy/ai.opencode.server.xenon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.opencode.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/opencode</string>
        <string>serve</string>
        <string>--mdns</string>
        <string>--port</string>
        <string>4096</string>
        <string>--mdns-domain</string>
        <string>xenon.opencode</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/ai/.local/share/opencode/server.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ai/.local/share/opencode/server-error.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/ai</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
