# deploy/sync-opencode-config.timer
[Unit]
Description=Hourly opencode config sync

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
