# deploy/sync-opencode-config.service
[Unit]
Description=Sync opencode config from git

[Service]
Type=oneshot
ExecStart=%h/.local/bin/sync-opencode-config
